# deploy/aluminatai-agent.service
[Unit]
Description=AluminatAI GPU Energy Monitoring Agent
Documentation=https://aluminatiai.com/docs/agent
After=network-online.target nvidia-persistenced.service
Wants=network-online.target
# Restart if NVML driver is temporarily unavailable (e.g. driver reload)
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=aluminatai
Group=aluminatai
EnvironmentFile=/etc/aluminatai/agent.env
Environment=DATA_DIR=/var/lib/aluminatai
Environment=LOG_DIR=/var/log/aluminatai

ExecStart=/usr/local/bin/aluminatiai
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=10s

# Give the agent time to flush on SIGTERM before SIGKILL
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aluminatai-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/aluminatai /opt/aluminatai-agent/data

[Install]
WantedBy=multi-user.target
